# rootfs/etc/services.d/scanservjs/run
#!/usr/bin/with-contenv sh
set -xve

cd "$APP_DIR"
exec node server.js
